# .config/systemd/user/squat_reminder.service
[Unit]
Description=Squat Reminder Notification

[Service]
ExecStart=/home/hector/dotfiles/scripts/squat_reminder.sh
